# nordvpn-reconnect.service
[Unit]
Description=NordVPN reconnector daemon

[Service]
ExecStart=/usr/bin/nordvpn-reconnector

[Install]
WantedBy=multi-user.target
